# filesystem/var/lib/update-rc.d/rc.local
{
  "MIME": "text/plain",
  "inode_type": "REG",
  "magic_type": "ASCII text",
  "mode": "-rw-r--r--",
  "sha1": "976f6eb3e377ede975e7d5fff9ae177ea18c2cea"
}