# filesystem/etc/init.d/rc.local
{
  "MIME": "text/x-shellscript",
  "inode_type": "REG",
  "magic_type": "POSIX shell script, ASCII text executable",
  "mode": "-rwxr-xr-x",
  "sha1": "ec77ee7a36127da7c56de28d6e4b7442a70e2877"
}