# filesystem/etc/init.d/rc.local
{
  "MIME": "text/x-shellscript",
  "inode_type": "REG",
  "magic_type": "POSIX shell script, ASCII text executable",
  "mode": "-rwxr-xr-x",
  "sha1": "e612460a1a3df8c31ed6a73e5987cd576cc8d2d2"
}