# filesystem/etc/network/if-up.d/000resolvconf
{
  "MIME": "text/x-shellscript",
  "inode_type": "REG",
  "magic_type": "POSIX shell script, ASCII text executable",
  "mode": "-rwxr-xr-x",
  "sha1": "a7ca1aee5980ec3fb14054cf00d243b36bcb3283"
}