# filesystem/etc/network/if-up.d/000resolvconf
{
  "MIME": "text/x-shellscript",
  "inode_type": "REG",
  "magic_type": "POSIX shell script, ASCII text executable",
  "mode": "-rwxr-xr-x",
  "sha1": "e333278f6b198813c0ec6cedb3b11e6e1ff4e840"
}