# filesystem/etc/network/if-up.d/000resolvconf
{
  "MIME": "text/x-shellscript",
  "inode_type": "REG",
  "magic_type": "POSIX shell script, ASCII text executable",
  "mode": "-rwxr-xr-x",
  "sha1": "0da9240caaef054ce9a596e31c949138d4964cf2"
}