# filesystem/etc/network/if-up.d/ntpdate
{
  "MIME": "text/x-shellscript",
  "inode_type": "REG",
  "magic_type": "POSIX shell script, ASCII text executable",
  "mode": "-rwxr-xr-x",
  "sha1": "9defefbbeb4317019b4d80e78631cac17b97de51"
}