# filesystem/etc/network/if-up.d/openssh-server
{
  "MIME": "text/x-shellscript",
  "inode_type": "REG",
  "magic_type": "POSIX shell script, ASCII text executable",
  "mode": "-rwxr-xr-x",
  "sha1": "66862f364f9519ea5d7ea3d7bf9d2fa3d9e73b0d"
}