# filesystem/etc/network/if-up.d/openssh-server
{
  "MIME": "text/x-shellscript",
  "inode_type": "REG",
  "magic_type": "POSIX shell script, ASCII text executable",
  "mode": "-rwxr-xr-x",
  "sha1": "42c49b826072d85e7330ba991683abe771cdb31f"
}