# filesystem/etc/network/if-up.d/openssh-server
{
  "MIME": "text/x-shellscript",
  "inode_type": "REG",
  "magic_type": "POSIX shell script, ASCII text executable",
  "mode": "-rwxr-xr-x",
  "sha1": "54f69018b149129233da8406bbd9b114ac24aa32"
}