# filesystem/etc/network/if-up.d/openssh-server
{
  "MIME": "text/x-shellscript",
  "inode_type": "REG",
  "magic_type": "POSIX shell script, ASCII text executable",
  "mode": "-rwxr-xr-x",
  "sha1": "84f15f60f04cb198777017cf21768c3682fab379"
}