# filesystem/etc/network/if-up.d/upstart
{
  "MIME": "text/x-shellscript",
  "inode_type": "REG",
  "magic_type": "POSIX shell script, ASCII text executable",
  "mode": "-rwxr-xr-x",
  "sha1": "fc47c64da7de7f19b3ece8ed3ca547f57fb1c51a"
}