# filesystem/etc/profile.d/bash_completion.sh
{
  "MIME": "text/plain",
  "inode_type": "REG",
  "magic_type": "ASCII text",
  "mode": "-rw-r--r--",
  "sha1": "b1e6c6abda374b993b1abb56ebfee74d95ebaada"
}