# filesystem/etc/profile.d/bash_completion.sh
{
  "MIME": "text/plain",
  "inode_type": "REG",
  "magic_type": "ASCII text",
  "mode": "-rw-r--r--",
  "sha1": "6f8eba81acd31ec3b556e70ed7ae7eff0709bbe5"
}