# filesystem/etc/rc.local
{
  "MIME": "text/x-shellscript",
  "inode_type": "REG",
  "magic_type": "POSIX shell script, ASCII text executable",
  "mode": "-rwxr-xr-x",
  "sha1": "a9ca22d71797c9bb824e1c9885f3412df5432cf2"
}