# filesystem/etc/update-motd.d/00-header
{
  "MIME": "text/x-shellscript",
  "inode_type": "REG",
  "magic_type": "POSIX shell script, ASCII text executable",
  "mode": "-rwxr-xr-x",
  "sha1": "e20c30dc089afb67d9ee20e7218ac92336544fd2"
}