# filesystem/etc/update-motd.d/10-help-text
{
  "MIME": "text/x-shellscript",
  "inode_type": "REG",
  "magic_type": "POSIX shell script, ASCII text executable",
  "mode": "-rwxr-xr-x",
  "sha1": "a3800960f849c60483b23fe9b7ea1eaba88db605"
}