# filesystem/etc/update-motd.d/10-help-text
{
  "MIME": "text/x-shellscript",
  "inode_type": "REG",
  "magic_type": "POSIX shell script, ASCII text executable",
  "mode": "-rwxr-xr-x",
  "sha1": "508b6639232e387e929da8212c64bccc6f85c11a"
}