# filesystem/etc/update-motd.d/50-motd-news
{
  "MIME": "text/x-shellscript",
  "inode_type": "REG",
  "magic_type": "POSIX shell script, ASCII text executable",
  "mode": "-rwxr-xr-x",
  "sha1": "0c9ee22c4207b046462e2ca90422a67f915f8944"
}