# filesystem/etc/update-motd.d/80-esm
{
  "MIME": "text/x-shellscript",
  "inode_type": "REG",
  "magic_type": "POSIX shell script, UTF-8 Unicode text executable",
  "mode": "-rwxr-xr-x",
  "sha1": "8dae8392570cfcabc57421da73f3ff9b0a4b5023"
}