# filesystem/etc/update-motd.d/80-livepatch
{
  "MIME": "text/x-shellscript",
  "inode_type": "REG",
  "magic_type": "POSIX shell script, ASCII text executable",
  "mode": "-rwxr-xr-x",
  "sha1": "97bbeccd813514cdde2cf13af818163e75621db4"
}