# filesystem/etc/update-motd.d/91-release-upgrade
{
  "MIME": "text/x-shellscript",
  "inode_type": "REG",
  "magic_type": "POSIX shell script, ASCII text executable",
  "mode": "-rwxr-xr-x",
  "sha1": "01e40d92804a7861b0f8aac24c3eb909a9862a96"
}