# filesystem/etc/update-motd.d/91-release-upgrade
{
  "MIME": "text/x-shellscript",
  "inode_type": "REG",
  "magic_type": "POSIX shell script, ASCII text executable",
  "mode": "-rwxr-xr-x",
  "sha1": "54353896cee3d8e88f00d89dbf2e6e4991b15112"
}