# filesystem/etc/update-motd.d/99-footer
{
  "MIME": "text/x-shellscript",
  "inode_type": "REG",
  "magic_type": "POSIX shell script, ASCII text executable",
  "mode": "-rwxr-xr-x",
  "sha1": "2ed1df7b71e5c2e0cb1a8313c764bda78d2945c9"
}